# deploy/systemd/clawde-reflect.service
[Unit]
Description=Clawde reflection job (extrai lições de janela de events/observations)
After=network.target

[Service]
Type=oneshot
WorkingDirectory=%h/.clawde
EnvironmentFile=-%h/.clawde/config/clawde.env
# P3.4 (T-114): comando dedicado `clawde reflect` lê events + observations
# recentes, monta prompt estruturado conforme contrato do reflector/AGENT.md
# e enfileira como task LOW com dedupKey horário.
ExecStart=%h/.clawde/dist/clawde reflect --since 24h

# Hardening Nivel 1 (ADR 0005).
PrivateTmp=yes
ProtectHome=read-only
ProtectSystem=strict
NoNewPrivileges=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources
CapabilityBoundingSet=
AmbientCapabilities=

ReadWritePaths=%h/.clawde

[Install]
WantedBy=default.target
