# deploy/systemd/clawde-restore-drill.service
[Unit]
Description=Clawde restore drill mensal

[Service]
Type=oneshot
WorkingDirectory=%h/.clawde
EnvironmentFile=-%h/.clawde/config/clawde.env
ExecStart=/bin/bash -lc "%h/.clawde/scripts/restore-drill.sh || { %h/.clawde/dist/clawde smoke-test --db /nonexistent/clawde-restore-drill-alert.db --output json >/dev/null 2>&1 || true; exit 1; }"

# Hardening idêntico ao worker (BEST_PRACTICES §10.4)
PrivateTmp=yes
ProtectHome=read-only
ProtectSystem=strict
NoNewPrivileges=yes
PrivateDevices=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources
CapabilityBoundingSet=
AmbientCapabilities=

ReadWritePaths=%h/.clawde /tmp

[Install]
WantedBy=default.target
