# deploy/systemd/clawde-worker.path
[Unit]
Description=Clawde worker fallback trigger (queue.signal)

[Path]
PathChanged=%h/.clawde/run/queue.signal
Unit=clawde-worker.service

[Install]
WantedBy=default.target
