# deploy/systemd/clawde-backup-hourly.timer
[Unit]
Description=Clawde backup snapshot agendado (hourly)

[Timer]
OnCalendar=hourly
Unit=clawde-backup-hourly.service
Persistent=true

[Install]
WantedBy=timers.target
